# system_files/kinoite/etc/profile.d/fcitx5-kde.sh
#!/bin/sh

[ "${GTK_IM_MODULE}" = "fcitx" ] && unset GTK_IM_MODULE
[ "${QT_IM_MODULE}" = "fcitx" ] && unset QT_IM_MODULE
