# system_files/usr/lib/systemd/system/docker-prune.timer
[Unit]
Description=Weekly Docker Prune

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
