# system_files/usr/lib/systemd/system/nix-gc.timer
[Unit]
Description=Weekly Nix Garbage Collection

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
