# system_files/shared/etc/profile.d/fastfetch.sh
#!/usr/bin/env bash

alias fastfetch="/usr/bin/fastfetch -c /usr/share/ublue-os/ublue-os.jsonc"
